# apps/api-go/deploy/systemd/nlp-api.service
[Unit]
Description=Network Learning Platform API
After=network.target mysql.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/network-learning-platform/apps/api-go
EnvironmentFile=/etc/nlp/api-go.env
ExecStart=/usr/local/go/bin/go run ./cmd/server
Restart=always
RestartSec=5
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal

# Basic hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/network-learning-platform

[Install]
WantedBy=multi-user.target
